# ustozlarga.conf
[program:ustozlargabot]
directory=/home/ubuntu/ustozlargabot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/ustozlargabot/logfile_err.log
stdout_logfile=/home/ubuntu/uzhemisbot/logfile.err